# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:haproxy]
command=/opt/haproxy/sbin/haproxy -f /opt/haproxy/conf/haproxy.cfg -p /tmp/haproxy.pid
autorestart=true
